# _scripts/im.ajin.stats.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- launchd agent: run update_stats_from_journal.py daily at 09:00 local.
     Install:  cp this plist → ~/Library/LaunchAgents/  and  launchctl load ~/Library/LaunchAgents/im.ajin.stats.plist
     Reload:   launchctl unload ... && launchctl load ...
     Logs:     ~/Library/Logs/mobetter-stats.{out,err}.log -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>im.ajin.stats</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/ajin/Documents/New project/personal/mobetter20.github.io/_scripts/update_stats_from_journal.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ajin/Documents/New project/personal/mobetter20.github.io</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/ajin/Library/Logs/mobetter-stats.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ajin/Library/Logs/mobetter-stats.err.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
